# systemd/home-gateway/wg-secondary.service
[Unit]
Description=Secondary WireGuard tunnel to fallback VPS relay
After=network-online.target wg-primary.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/wg-quick up /etc/wireguard/wg-fallback.conf
ExecStop=/usr/bin/wg-quick down /etc/wireguard/wg-fallback.conf
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
